# root/etc/s6-overlay/s6-rc.d/unifi-controller-fixattrs/run
#!/command/with-contenv sh

chown -R appuser:appuser /config /run/unifi /var/run/unifi /usr/lib/unifi/run
chmod -R g+w /config /run/unifi /var/run/unifi /usr/lib/unifi/run
