# root/etc/s6-overlay/s6-rc.d/unifi-controller/run
#!/command/with-contenv sh

if [ -z "${XMX}" ]; then
    XMX="1024"
fi
if [ -z "${XMS}" ]; then
    XMS="256"
fi

exec s6-setuidgid appuser java -Xms"${XMS}M" -Xmx"${XMX}M" -Dlog4j2.formatMsgNoLookups=true -jar /usr/lib/unifi/lib/ace.jar start
